# systemd/kiosk.service
[Unit]
Description=DashPi Chromium Kiosk
After=network.target

[Service]
User=dashpi
Environment=DISPLAY=:0
ExecStart=/usr/local/dashpi/scripts/kiosk.sh
Restart=always
RestartSec=5

[Install]
WantedBy=graphical.target