# systemd/wifi-watchdog.service
[Unit]
Description=DashPi WiFi Watchdog
After=network.target

[Service]
User=dashpi
ExecStart=/usr/local/dashpi/scripts/wifi-check.sh
Restart=always
RestartSec=300

[Install]
WantedBy=multi-user.target